# solar_grabber.service
[Unit]
Description=Solar Grabber

[Service]
ExecStart=/bin/sun-status-grabber

[Install]
WantedBy=multi-user.target

